# meteo.crontab
*/5  *  * * *   root    /home/pi/rpi-meteo/sensors.sh
*/15  *    * * *   root    /home/pi/rpi-meteo/weather.sh
